# playbooks/files/k3s-agent.service.d/10.service.conf
# /etc/systemd/system/k3s.service.d/10.service.conf 
[Service]
Environment=GOGC=30
MemoryLimit=768M